# scripts/ethnomusicology-frontend.service
[Unit]
Description=Ethnomusicology Next.js Frontend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ethnomusicology/frontend-current
ExecStart=/home/ubuntu/.bun/bin/bun --bun next start -p 3000
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
